# services/logrotate-wireguard.timer
[Unit]
Description=Run logrotate-wireguard every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target