# services/logrotate-wireguard.service
[Unit]
Description=Rotate WireGuard logs

[Service]
Type=oneshot
ExecStart=/usr/local/bin/logrotate-wireguard.sh