# services/wg_systemd.service
[Unit]
Description=WireGuard systemd logger

[Service]
ExecStart=/bin/sh -c 'journalctl -f | grep wireguard >> /etc/wireguard/logs/wg_systemd.log'
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target