# key_management_scripts/wg-sync-key.service
[Unit]
Description=WireGuard Key Sync

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /usr/local/bin/wg-key-sync.py
